# s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv bash

echo "⚙️ Starting worker..."

cd /app
export DATABASE_PATH=/app/app_data/data/matchexec.db
exec s6-setuidgid abc tsx processes/worker/index.ts