# s6-overlay/s6-rc.d/matchexec-web/run
#!/command/with-contenv bash

echo "Starting MatchExec web server..."

cd /app
exec s6-setuidgid abc node server.js