# s6-overlay/s6-rc.d/scheduler/run
#!/command/with-contenv bash

echo "⏰ Starting scheduler..."

cd /app
exec s6-setuidgid abc tsx processes/scheduler/index.ts