# s6-overlay/s6-rc.d/scheduler/run
#!/command/with-contenv bash

echo "Starting scheduler..."

cd /app
exec s6-setuidgid abc node dist/scheduler.js