# s6-overlay/s6-rc.d/discord-bot/run
#!/command/with-contenv bash

echo "Starting Discord bot..."

cd /app
exec s6-setuidgid abc tsx processes/discord-bot/index.ts