# s6-overlay/s6-rc.d/discord-bot/run
#!/command/with-contenv bash

echo "Starting Discord bot..."

cd /app
exec s6-setuidgid abc node dist/discord-bot.js