# s6-overlay/s6-rc.d/db-migrator/up
s6-setuidgid abc /bin/bash /app/scripts/run-migrations.sh