# s6-overlay/s6-rc.d/db-migrator/up
#!/command/with-contenv bash

echo "📊 Running database migrations and seeding..."

# Run migrations as the abc user (blocking until complete)
cd /app
s6-setuidgid abc tsx scripts/migrate-background.ts

echo "✅ Database initialization completed"